# mconnect/.config/runit/sv/mconnect/run
#!/bin/sh
set -efu
export DISPLAY=:0.0
[ -r conf ] && ./conf
exec mconnect ${OPTS:-}
