# newsboat/.config/systemd/user/newsboat.timer
[Unit]
Description=Newsboat Feed Updater
After=network-online.target

[Timer]
OnActiveSec=2m
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
