# mpd/.config/systemd/user/mpd.service.d/10-create-dirs.conf
[Service]
ExecStartPre=/usr/bin/mkdir -p %C/mpd %h/.local/share/mpd
