# kdeconnect/.config/runit/sv/kdeconnect/run
#!/bin/sh -efu
export DISPLAY=${DISPLAY:-:0.0}
exec /usr/lib/kdeconnectd
