# newsboat/.config/runit/sv/newsboat/run
#!/bin/sh
exec snooze -H/1 newsboat -x reload
