# newsboat/.config/runit/sv/newsboat/run
#!/bin/sh
[ -r conf ] && . ./conf
exec snooze -H/1 newsboat -x reload
